# fleet-units/global.service
[Unit]
Description=Global Test

[Service]
ExecStart=/usr/bin/true

[X-Fleet]
Global=true
